# 1GB-Server-aws/deploy/1gb-server.service
[Unit]
Description=1GB Server Flask App (Gunicorn)
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/1gb-server

Environment="PATH=/home/ubuntu/1gb-server/venv/bin"
Environment="PYTHONUNBUFFERED=1"

ExecStart=/home/ubuntu/1gb-server/venv/bin/gunicorn -c deploy/gunicorn.conf.py main:app
Restart=always
RestartSec=3

# Hardening (safe defaults)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
