# 1GB-Server-aws/app/deploy/1gb-server.service
[Unit]
Description=1GB Server (Gunicorn) - Production Monitoring Dashboard
After=network.target

[Service]
Type=simple

# Change these paths to match where you clone the repo
WorkingDirectory=/opt/1gb-server

Environment=APP_NAME=1GB%20Server
Environment=APP_VERSION=2.0.0
Environment=DISK_PATH=/
Environment=MAX_APP_MEMORY_MB=1024
Environment=APP_MEMORY_WARN_RATIO=0.85
Environment=WARMUP_SECONDS=2

# Gunicorn bind is localhost:8000; Nginx listens on :80
Environment=PORT=8000
Environment=BIND=127.0.0.1:8000

# t3.micro-friendly defaults
Environment=WEB_CONCURRENCY=1
Environment=GUNICORN_THREADS=2
Environment=LOG_LEVEL=INFO

ExecStart=/opt/1gb-server/venv/bin/gunicorn -c /opt/1gb-server/gunicorn.conf.py main:app
Restart=always
RestartSec=2

User=ec2-user
Group=ec2-user

# Logs go to journald so you can view them from any new PuTTY session:
#   sudo journalctl -u 1gb-server -f
StandardOutput=journal
StandardError=journal

# Slightly safer sandboxing while still allowing psutil reads
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
